# .config/autostart/Xilinx Information Center.desktop
[Desktop Entry]
Encoding=UTF-8
Type=Application
Name=Xilinx Information Center
Comment=Xilinx Information Center
Icon=/home/foo/Xilinx/xic/data/images/xlnx_logo.png
Exec=/home/foo/Xilinx/xic/xic -c -t
